# .cmake/PROCONT-EX.service.in
[Unit]
Description=Procont Execution Environment (Local)
After=network.target

[Service]
Type=simple
User=${USER}
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/${PROJECT_NAME}
WorkingDirectory=${CMAKE_INSTALL_PREFIX}
Environment="CONFIG_DIR=${CMAKE_INSTALL_PREFIX}/etc"
Environment="LOG_DIR=${CMAKE_INSTALL_PREFIX}/logs"
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
